# deploy/systemd/shophosting-system-backup.service
[Unit]
Description=ShopHosting System Backup
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/shophosting/scripts/system-backup.sh
StandardOutput=journal
StandardError=journal
